# mixed-criticality-orchestrator/storage/apps/asil/b/acc.container
[Container]
Image=localhost/vehicle-detection:latest
ContainerName=acc
PublishPort=5000:5000
Network=host
#Exec=sleep infinity