# mixed-criticality-orchestrator/storage/apps/asil/b/laneassist.container
[Container]
Image=localhost/lane-assistant:latest
ContainerName=laneassist
PublishPort=5009:5009
Network=host
#Exec=sleep infinity