# mixed-criticality-orchestrator/storage/apps/asil/b/rearcamera.container
[Container]
Image=localhost/rear-camera-app
ContainerName=rearcamera
PublishPort=8025:8025
Network=host