# mixed-criticality-orchestrator/FFI/CPU/usecase3/qm-app3.container
[Container]
ContainerName=Cpu-stress-app3
#Image=localhost/cpu-stress-app
Image=localhost/qm-cpu-app
Network=host
[Service]
CPUWeight=2048

