# mixed-criticality-orchestrator/storage/apps/qm/infotainment/games/pacman.container
[Container]
Image=localhost/pacman-game:latest
ContainerName=pacman
Volume=/tmp/apache-logs:/usr/local/apache2/logs
PublishPort=8022:80
Network=host
