# mixed-criticality-orchestrator/storage/apps/qm/infotainment/ott/netflix.container
[Container]
Image=localhost/netflix
ContainerName=netflix
PublishPort=8024:8024
Network=host