# mixed-criticality-orchestrator/FFI/PID_LIMITS/acc.container
[Container]
Image=localhost/vehicle-detection-frame:latest
ContainerName=acc
PublishPort=5000:5000
Network=host
